# rootfs/static/var/spool/cron/crontabs/root
00   00 * * * /usr/sbin/ntpdate 1.kr.pool.ntp.org && /sbin/hwclock -u -w
00-59 * * * * check_taskmon &
